# service/cloudflare-dns.service
[Unit]
Description=cloudflare-dns

[Service]
ExecStart=docker run --rm go 
