# service/cloudflare-dns.service
[Unit]
Description=cloudflare-dns

[Service]
ExecStart=docker run --rm u6du/cloudflare-dns
